# graveyard/glustermount
#!/sbin/openrc-run

description="Mounts gluster volumes"

depend()
{
        local opts mywant=""
        after root
        want $mywant
        need glusterd
        before registry
        before lighttpd
        use dns
        use root
}

start()
{
        ebegin "Mounting gluster volume"
        einfo "gluster volumes:"$(gluster volume list)
        for v in $(gluster volume list); do 
           einfo "mounting localhost:$v to $shared_mounts/$v"
	   [ ! -d $shared_mounts/$v ] && mkdir "/data/$v"
           mount -t glusterfs localhost:$v $shared_mounts/$v
        done
        rc=0
        eend "Done"
        return $rc
}

stop()
{
  ebegin "Unmounting gluster volumes"
        for v in $(gluster volume list); do 
           umount $shared_mounts/$v
        done
        rc=0
        eend "Done"
        return $rc 
}
