# ansible/roles/gentoo/files/init.d/set_cpu_performance
#!/sbin/openrc-run

description="CPU performance settings"

depend()
{
        local opts mywant=""
        after root
}

start()
{
ebegin "starting"
LOGFILE="/tmp/set_cpu"
echo "date:$(date)" > "${LOGFILE}"
cpupower frequency-set -g ondemand

eend "succeeded"
return 0
}

stop()
{
  true
}
