# roles/gentoo/files/etc/init.d/brickmount
#!/sbin/openrc-run

description="Mounts bricks for gluster volumes"
glusterconf=/var/lib/glusterd
my_ips=$(cat /proc/net/fib_trie | awk '/32 host/ { print f } {f=$2}')
gluster_dir=/data/gluster
brick_disk=${gluster_dir}/disks/

depend()
{
        local opts mywant=""
        after root
        want $mywant
        use dns
        use root
}

action_on_bricks()
{
  action="$1"
  ebegin "Mounting bricks"
  for v in $(find ${glusterconf}/vols/ -maxdepth 1 -mindepth 1 -type d -printf "%f\n"); do
    for loc in $(grep -l "^hostname=$(hostname)" ${glusterconf}/vols/${v}/bricks/* );do
      echo "found volume $v and brick $loc"
      mp=$(grep "^path=" "${loc}" | sed -e 's/^path=//g')
      app_name=$(dirname "${mp}" | xargs basename)
      vol_name=$(basename "${mp}")
      disk="${brick_disk}/app:${app_name}-vol:${vol_name}"

      echo "${action} ${disk} ${mp} (app_name=${app_name} vol=${vol_name}"
      if [ "${action}" = "mount" ]; then
        mount ${disk} ${mp}
      elif [ "${action}" = "umount" ]; then
        umount ${mp}
      fi
    done
  done
  rc=0
  eend "Done"
  return $rc
}

start()
{
  action_on_bricks mount
}

stop()
{
  action_on_bricks umount
}
