# roles/gentoo/files/etc/init.d/registry
#!/sbin/openrc-run

description="run registry"

depend()
{
        local opts mywant=""
        need glusterd
        want $mywant
        use dns
        use root
}

load_image(){
    registry_content="/shared/registry/docker/registry/v2"
    cd $(mktemp --directory)
    reg_sha=$(cut -d : -f 2 ${registry_content}/repositories/local/registry/_manifests/tags/${service_image##*:}/current/link )
    config_sha=$(cat ${registry_content}/blobs/sha256/$(echo "$reg_sha"| cut -c -2)/$reg_sha/data | jq -r .config.digest  | cut -d : -f 2)
    echo '[{"Config": "'$config_sha'.json","RepoTags": ["'${service_image}'"],"Layers": []}]' > manifest.json
    ln -s  ${registry_content}/blobs/sha256/$(echo "$config_sha" | cut -c -2)/$config_sha/data $config_sha.json
    tars=""; for blob_sha in $(cat ${registry_content}/blobs/sha256/$(echo "$reg_sha"| cut -c -2)/$reg_sha/data | jq -r .layers[].digest | cut -d : -f 2); do
      ln -s ${registry_content}/blobs/sha256/$(echo "$blob_sha" | cut -c -2)/$blob_sha/data $blob_sha.tar; 
      nm=$(jq '.[0].Layers += ["'$blob_sha.tar'"]' manifest.json ); 
      echo "$nm" > manifest.json; tars="$tars $blob_sha.tar"; 
    done
    tar  --dereference -cf reg.tar  manifest.json $config_sha.json $tars
    podman load -i reg.tar
    cd -
}

start()
{
        ebegin "Starting registry"
	einfo "registry: ${service_hostname} host_volume: ${host_volume=} image:${service_image}"
	hn=$(echo "${service_hostname}" | cut -d : -f 1)

        ip=$(getent hosts "${hn}" | awk '{print $1}')
        ip addr | grep -q $ip
        #(ping -q -c 1 "${service_hostname}" 2>&1) > /dev/null
        ret_ip=$?
        if [ $ret_ip -ne 0 ];then 
          einfo "ip of host ${hn} not assigned here, nothing to do!"
          eend "Done"
          return 0
        fi

        # ensure image exists
	if ! (podman image ls | awk '{print $1":"$2}' | grep -q "^${service_image}" ); then
          load_image
	fi

        # claim IP
        #ip address add "${ipconf}" dev wlan0
        podman rm -f registry
        podman run --rm --privileged -d --name registry -p 1443:443 -v ${host_volume}:/var/lib/registry \
                -e REGISTRY_STORAGE_DELETE_ENABLED=true \
                -e REGISTRY_HTTP_ADDR=0.0.0.0:443 \
                -e REGISTRY_HTTP_TLS_CERTIFICATE=/var/lib/registry/certs/myregistry.adm13.crt \
                -e REGISTRY_HTTP_TLS_KEY=/var/lib/registry/certs/myregistry.adm13.key \
        "${service_image}" serve /var/lib/registry/config.yml

        rc=0
        eend "Done"
        return $rc
}

stop()
{
  ebegin "Stopping registry"
  podman stop registry
  rc=0
  return $rc
}
